# docs/tldraw-ai-board.service
[Unit]
Description=AI Board Fastify Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/tldraw-ai-board
Environment=NODE_ENV=production
Environment=PORT=3333
EnvironmentFile=/srv/tldraw-ai-board/.env
ExecStart=/usr/local/bin/node dist/server/server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
